# library_c/unicornd/unicornd.service
[Unit]
Description=Unicorn HAT daemon

[Service]
ExecStart=/usr/sbin/unicornd

[Install]
WantedBy=multi-user.target
